# coreos/etc/systemd/system/prune-docker-images.timer
[Unit]
Description=Run prune-docker-images.service once a day.

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
